# backen/scraper/com.prediactor.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prediactor.scraper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/daniel_mba/Documents/GitHub/Prediactor/backen/scraper/run_daily.sh</string>
  </array>

  <!-- Run every day at 02:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Don't run the moment it's loaded; only at 2 AM -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/daniel_mba/Documents/GitHub/Prediactor/backen/scraper/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/daniel_mba/Documents/GitHub/Prediactor/backen/scraper/logs/launchd.err.log</string>
</dict>
</plist>
